# walker/.config/systemd/user/walker.service
[Unit]
Description=Walker Launcher Service
After=elephant.service

[Service]
Type=simple
ExecStart=/usr/bin/walker --gapplication-service
Restart=on-failure

[Install]
WantedBy=graphical-session.target
